# playbook/roles/vault/templates/vault.service.j2
[Unit]
Description=Vault process

[Service]
User=root
Restart=on-failure
ExecStart=/usr/local/bin/vault server -config={{vault.configBase}}/config.hcl

[Install]
WantedBy=multi-user.target